# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:rabbitmq]
command=/usr/sbin/rabbitmq-server
autostart=true
autorestart=true
priority=1
stdout_logfile=/var/log/supervisor/rabbitmq.log
stderr_logfile=/var/log/supervisor/rabbitmq_err.log
environment=HOME="/var/lib/rabbitmq"

[program:rabbitmq-setup]
command=/bin/bash -c "sleep 10 && rabbitmqctl add_user rlvr rlvr_password || true && rabbitmqctl set_permissions -p / rlvr '.*' '.*' '.*' || true"
autostart=true
autorestart=false
startsecs=0
priority=2
stdout_logfile=/var/log/supervisor/rabbitmq-setup.log
stderr_logfile=/var/log/supervisor/rabbitmq-setup_err.log

[program:qdrant]
command=/usr/local/bin/qdrant
directory=/workspace/qdrant_storage
autostart=true
autorestart=true
priority=3
stdout_logfile=/var/log/supervisor/qdrant.log
stderr_logfile=/var/log/supervisor/qdrant_err.log

[program:ollama]
command=/bin/bash -c "ollama serve"
autostart=true
autorestart=true
priority=4
stdout_logfile=/var/log/supervisor/ollama.log
stderr_logfile=/var/log/supervisor/ollama_err.log
environment=OLLAMA_HOST="0.0.0.0:11434",OLLAMA_GPU_LAYERS="999",CUDA_VISIBLE_DEVICES="0"

[program:ollama-pull-model]
command=/bin/bash -c "sleep 30 && ollama pull llama3.2:3b"
autostart=true
autorestart=false
startsecs=0
priority=5
stdout_logfile=/var/log/supervisor/ollama-pull.log
stderr_logfile=/var/log/supervisor/ollama-pull_err.log

[program:document-ingestion]
command=/bin/bash -c "sleep 40 && cd /workspace/rlvr-pdf-chat/services/document-ingestion && exec uvicorn src.main:app --host 0.0.0.0 --port 8002"
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/document-ingestion.log
stderr_logfile=/var/log/supervisor/document-ingestion_err.log
environment=PYTHONPATH="/workspace/rlvr-pdf-chat/services/document-ingestion:/workspace/rlvr-pdf-chat",QDRANT_HOST="localhost",QDRANT_PORT="6333",RABBITMQ_URL="amqp://rlvr:rlvr_password@localhost:5672/"

[program:qa-orchestrator]
command=/bin/bash -c "sleep 45 && cd /workspace/rlvr-pdf-chat/services/qa-orchestrator && exec uvicorn src.main:app --host 0.0.0.0 --port 8001"
autostart=true
autorestart=true
priority=11
stdout_logfile=/var/log/supervisor/qa-orchestrator.log
stderr_logfile=/var/log/supervisor/qa-orchestrator_err.log
environment=PYTHONPATH="/workspace/rlvr-pdf-chat/services/qa-orchestrator:/workspace/rlvr-pdf-chat",QDRANT_URL="http://localhost:6333",OLLAMA_URL="http://localhost:11434",RABBITMQ_URL="amqp://rlvr:rlvr_password@localhost:5672/"

[program:api-gateway]
command=/bin/bash -c "sleep 50 && cd /workspace/rlvr-pdf-chat/services/api-gateway && exec uvicorn src.main:app --host 0.0.0.0 --port 8000"
autostart=true
autorestart=true
priority=12
stdout_logfile=/var/log/supervisor/api-gateway.log
stderr_logfile=/var/log/supervisor/api-gateway_err.log
environment=PYTHONPATH="/workspace/rlvr-pdf-chat/services/api-gateway:/workspace/rlvr-pdf-chat",QA_ORCHESTRATOR_URL="http://localhost:8001",DOC_INGESTION_SERVICE_URL="http://localhost:8002"

[program:verification-worker]
command=/bin/bash -c "sleep 55 && cd /workspace/rlvr-pdf-chat/workers/verification-worker && exec python -m src.worker"
autostart=true
autorestart=true
priority=13
stdout_logfile=/var/log/supervisor/verification-worker.log
stderr_logfile=/var/log/supervisor/verification-worker_err.log
environment=PYTHONPATH="/workspace/rlvr-pdf-chat/workers/verification-worker:/workspace/rlvr-pdf-chat",RABBITMQ_URL="amqp://rlvr:rlvr_password@localhost:5672/",OLLAMA_URL="http://localhost:11434",OLLAMA_MODEL="llama3.2:3b"

[program:dataset-worker]
command=/bin/bash -c "sleep 60 && cd /workspace/rlvr-pdf-chat/workers/dataset-generation-worker && exec python -m src.worker"
autostart=true
autorestart=true
priority=14
stdout_logfile=/var/log/supervisor/dataset-worker.log
stderr_logfile=/var/log/supervisor/dataset-worker_err.log
environment=PYTHONPATH="/workspace/rlvr-pdf-chat/workers/dataset-generation-worker:/workspace/rlvr-pdf-chat",RABBITMQ_URL="amqp://rlvr:rlvr_password@localhost:5672/",TRAINING_DATA_DIR="/workspace/rlvr-pdf-chat/data/training_data",DPO_DATA_DIR="/app/data/dpo_data"

[program:streamlit]
command=/bin/bash -c "sleep 65 && cd /workspace/rlvr-pdf-chat/ui/streamlit && exec streamlit run src/app_simple.py --server.port 8501 --server.address 0.0.0.0"
autostart=true
autorestart=true
priority=15
stdout_logfile=/var/log/supervisor/streamlit.log
stderr_logfile=/var/log/supervisor/streamlit_err.log
environment=PYTHONPATH="/workspace/rlvr-pdf-chat/ui/streamlit:/workspace/rlvr-pdf-chat",API_GATEWAY_URL="http://localhost:8000"

